# ansible/playbooks/roles/docker/templates/docker-binfmt-register.service.j2
[Unit]
Description=Register binfmt emulators via Docker
After=docker.service network-online.target
Wants=docker.service network-online.target
Requires=proc-sys-fs-binfmt_misc.mount
After=proc-sys-fs-binfmt_misc.mount

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sh -c '/usr/bin/mountpoint -q /proc/sys/fs/binfmt_misc || /usr/bin/mount -t binfmt_misc none /proc/sys/fs/binfmt_misc'
ExecStartPre=/bin/sh -c 'for attempt in 1 2 3 4 5 6 7 8 9 10; do /usr/bin/docker info >/dev/null 2>&1 && exit 0; /usr/bin/sleep 1; done; exit 1'
ExecStart=/usr/bin/docker run --privileged --rm {{ docker_binfmt_image }} --install all

[Install]
WantedBy=multi-user.target
